# config/docker/runit/mollypictures/run
#!/bin/bash

exec /app/rel/mollypictures/bin/mollypictures foreground
